# debian/linstor-proxmox.postinst
#!/bin/sh

invoke-rc.d pvedaemon restart || :
